# statemint-app.service
[Unit]
Description=StateMint Next.js Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/your/app  # Update this path
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

# Logging
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=statemint-app

[Install]
WantedBy=multi-user.target 